# linux/usr/lib/systemd/system/wiimotedev-io.target
[Unit]
Description=Wiiremote mouse/keyboard/joystick emulator