# linux/usr/lib/systemd/system/wiimotedev-daemon.target
[Unit]
Description=Wiiremote connection service and data provider