# linux/etc/init.d/gentoo/wiimotedev-io
#!/sbin/runscript
# Distributed under the terms of the GNU Lesser General Public License, v2.1 or later
# Copyright (C) 2010 Bartłomiej Burdukiewicz
# Contact: bartlomiej.burdukiewicz@gmail.com   

PIDFILE=/var/run/wiimotedev-io.pid
SERVICE=/usr/sbin/wiimotedev-io
GROUP=io
CONFIG=/etc/wiimotedev/io.conf

depend() {
	need dbus 
	after wiimotedev-daemon
}

start() {
  ebegin "Starting Wiimotedev-io"

  if [ -z $(grep ^$GROUP: /etc/group) ]; then
    groupadd $GROUP &> /dev/null

    if [ ! -n $(grep ^$GROUP: /etc/group) ]; then
      eend 1
      echo "Cannot add new group $GROUP"
      exit 1
    fi

  fi

  if [ -f "$PIDFILE" ]; then
    eend 1
    echo "Already running. Please remove $PIDFILE to force start."
    exit 1  
  fi
  
  if [ -x "$SERVICE" ]; then
    start-stop-daemon --start -q -x "$SERVICE" --group "$GROUP"
  fi

  eend 0
}

stop() {
  ebegin "Stopping Wiimotedev-io"

  if [ -f "$PIDFILE" ]; then
    start-stop-daemon --stop --quiet --pidfile "$PIDFILE"
    rm -f $PIDFILE &> /dev/null
  else
    eend 1
    echo "Please start service first."
    exit 1
  fi

  eend 0
}

restart() {
	stop
        sleep 1
	start
}


