# init.d/client/gentoo/wiimotedev-client
#!/sbin/runscript
# Distributed under the terms of the GNU Lesser General Public License, v2.1 or later
# Copyright (C) 2010 Bartłomiej Burdukiewicz
# Contact: dev.strikeu@gmail.com   

PIDFILE=/var/run/wiimotedev-client.pid
SERVICE=/usr/sbin/wiimotedev-client
GROUP=wiimotedev
CONFIG=/etc/wiimotedev/client.conf

depend() {
	need dbus 
	after wiimotedev-daemon
}

start() {
  ebegin "Starting Wiimotedev-client"

  if [ -z $(cat /etc/group | grep $GROUP) ]; then
    groupadd $GROUP &> /dev/null

    if [ ! -n $(cat /etc/group | grep $GROUP) ]; then
      eend 1
      echo "Cannot add new group $GROUP"
      exit 1
    fi

  fi

  if [ -f "$PIDFILE" ]; then
    eend 1
    echo "Already running. Please remove $PIDFILE to force start."
    exit 1  
  fi
  
  if [ -x "$SERVICE" ]; then
    start-stop-client --start -q -x "$SERVICE" --group "$GROUP"
  fi

  eend 0
}

stop() {
  ebegin "Stopping Wiimotedev-client"

  if [ -f "$PIDFILE" ]; then
    kill -s 3 $(cat $PIDFILE) &> /dev/null
    rm -f $PIDFILE &> /dev/null
  else
    eend 1
    echo "Please start service first."
    exit 1
  fi

  eend 0
}

restart() {
	stop
        sleep 1
	start
}


