# init.d/uinput/gentoo/wiimotedev-uinput
#!/sbin/runscript
# Distributed under the terms of the GNU Lesser General Public License, v2.1 or later
# Copyright (C) 2010 Bartłomiej Burdukiewicz
# Contact: dev.strikeu@gmail.com   

PIDFILE=/var/run/wiimotedev-uinput.pid
SERVICE=/usr/sbin/wiimotedev-uinput
GROUP=uinput
CONFIG=/etc/wiimotedev/uinput.conf

opts="reload"

# int wiimotedev_check_group()

function wiimotedev_check_group {
  if [ -z $(cat /etc/group | grep "$GROUP") ]; then
    echo " -- Adding new group \"${GROUP}\""
    groupadd $GROUP

    if [ -n $(cat /etc/group | grep "$GROUP") ]; then
      echo " -- Group \"${GROUP}\" added successful"
    else
      echo " -- Cannot add new group \"${GROUP}\""
      return 1
    fi
  fi
}

# int wiimotedev_start()

function wiimotedev_start() {
  wiimotedev_check_group

  if [ ! -x $SERVICE ]; then
    echo " -- ${SERVICE} file is not executable"
    return 1
  fi

  start-stop-daemon --start -q -x "$SERVICE" --group "$GROUP"
}

# int wiimotedev_stop()

function wiimotedev_stop() {
  wiimotedev_check_group

  if [ -f "$PIDFILE" ]; then
    start-stop-daemon --stop -q -p "$PIDFILE" --group "$GROUP"
  fi
}

depend() {
	need dbus 
}

start() {
	ebegin "Starting wiimotedev-client"
	wiimotedev_start
	eend $?
}

restart() {
	ebegin "Restarting wiimotedev-client"
	wiimotedev_start
	wiimotedev_stop
	eend $?
}

stop() {
	ebegin "Stopping wiimotedev-client"
	wiimotedev_stop
	eend $?
}