# initd/gentoo/wiimotedev
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License, v2 or later


depend() {
	need dbus bluetooth
}

start() {
	ebegin "Starting Wiimotedev-daemon"
        start-stop-daemon --start -q --exec /usr/sbin/wiimotedev
	eend $?
}

stop() {
	ebegin "Stopping Wiimotedev-daemon"
	start-stop-daemon --stop -q --pidfile /var/run/wiimotedev.pid
	eend $?
}
