# initd/gentoo/wiimotedev-daemon
#!/sbin/runscript
# Distributed under the terms of the GNU Lesser General Public License, v2.1 or later
# Copyright (C) 2008-2010 Bartłomiej Burdukiewicz
# Contact: dev.strikeu@gmail.com   

depend() {
	need dbus bluetooth
}

start() {
	ebegin "Starting Wiimotedev-daemon"
        start-stop-daemon --start -q --exec /usr/sbin/wiimotedev-daemon
	eend $?
}

stop() {
	ebegin "Stopping Wiimotedev-daemon"
	start-stop-daemon --stop -q --pidfile /var/run/wiimotedev-daemon.pid
	eend $?
}
